# packaging/systemd/helm-packd.service
[Unit]
Description=Helm C++ local pack service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=helm
Group=helm
EnvironmentFile=/etc/helm/helm-runtime.env
Environment=HELM_BIND=0.0.0.0
WorkingDirectory=/opt/helm
ExecStart=/opt/helm/bin/helm-packd 8091
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/etc/helm /var/lib/helm /var/cache/helm /var/log/helm /srv/helm
StandardOutput=append:/var/log/helm/helm-packd.log
StandardError=append:/var/log/helm/helm-packd.err.log

[Install]
WantedBy=multi-user.target
