# packaging/systemd/helm-basemap-cache.service
[Unit]
Description=Helm C++ basemap cache and online-fill service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=helm
Group=helm
EnvironmentFile=/etc/helm/helm-runtime.env
Environment=HELM_BIND=0.0.0.0
WorkingDirectory=/opt/helm
ExecStart=/opt/helm/bin/helm-basemap-cache 8095
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
ReadWritePaths=/etc/helm /var/lib/helm /var/cache/helm /var/log/helm /srv/helm
StandardOutput=append:/var/log/helm/helm-basemap-cache.log
StandardError=append:/var/log/helm/helm-basemap-cache.err.log

[Install]
WantedBy=multi-user.target
