# packaging/launchd/com.6thelement.helm-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.6thelement.helm-server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/helm/bin/helm-server</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/opt/helm</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HELM_BIND</key>
    <string>0.0.0.0</string>
    <key>HELM_PORT</key>
    <string>8080</string>
    <key>HELM_WEB_ROOT</key>
    <string>/opt/helm/web</string>
    <key>HELM_CONFIG</key>
    <string>/etc/helm</string>
    <key>HELM_RUNTIME_DIR</key>
    <string>/var/lib/helm/runtime</string>
    <key>HELM_S57_DATA</key>
    <string>/var/lib/helm/runtime/s57data</string>
    <key>HELM_TCDATA_DIR</key>
    <string>/var/lib/helm/runtime/tcdata</string>
    <key>HELM_SENC_DIR</key>
    <string>/var/cache/helm/senc</string>
    <key>HELM_TIDES_CACHE_DIR</key>
    <string>/var/cache/helm/tides</string>
    <key>HELM_USER_DATA_ROOT</key>
    <string>/var/lib/helm/data</string>
    <key>HELM_ENC</key>
    <string>/var/lib/helm/runtime/enc/US5FL4CR/US5FL4CR.000</string>
    <key>HELM_TILES_NO_WARMUP</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/var/log/helm/helm-server.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/helm/helm-server.err.log</string>
</dict>
</plist>
